# supervisord.conf
[supervisord]
nodaemon=true

[program:int]
command=/usr/local/bin/init init
